# contrib/systemd/smtpw.target
[Unit]
Description=SMTPw control all instances
